# rootfs-focal/etc/systemd/system/hostapd@.service.d/override.conf
[Unit]
Before=ssh-fix-reboot.service

[Service]
SyslogIdentifier=hostapd
ExecStart=
ExecStart=/usr/local/sbin/hostapd-launch --device %i

