# rootfs-focal/etc/systemd/system/hostapd@.service.d/override.conf
[Service]
SyslogIdentifier=hostapd
ExecStart=
ExecStart=/usr/local/sbin/hostapd-launch %i

